# deploy/whisp.service
[Unit]
Description=whisp transcription service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/whisp --config /opt/whisp/config.json
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
